# cron/cron.monthly/run_monthly
#!/bin/bash
# 月度缓存清理 — 智能清理退市/除权缓存，释放磁盘空间
# 默认模式：智能清理（退市+除权+损坏）
# --all 模式：删除所有缓存（下次运行将全量重新拉取）
SCRIPT_DIR=$(cd "$(dirname "$0")" && pwd)
PROJECT_ROOT=$(cd "$SCRIPT_DIR/../.." && pwd)
source "$PROJECT_ROOT/cron/_common.sh"

init_env
LOG_FILE=$LOG_DIR/monthly.log

if [ "$1" = "--all" ]; then
    log_info "开始全量清除缓存"
    rm -rf "$PROJECT_ROOT/quantia/cache/hist"/*
    log_info "已全量清除所有历史K线缓存"
else
    log_info "开始智能清理缓存"
    "$PYTHON_BIN" -c "
import sys; sys.path.insert(0, '$PROJECT_ROOT')
from quantia.core.stockfetch import clean_expired_cache
clean_expired_cache()
" >> "$LOG_FILE" 2>&1
    log_info "缓存清理完成"
fi

# ── 月度财务数据更新 ──
# 披露季后(5/9/11月)做全量采集，捕获新披露报告期与财报修订；
# 其余月份增量（按整只跳过已追平股票，分钟级完成），避免重复抓取。
cd "$PROJECT_ROOT"
MONTH=$((10#$(date +%m)))
case "$MONTH" in
    5|9|11) FIN_MODE="" ;;
    *)      FIN_MODE="--incremental" ;;
esac
if [ -z "$FIN_MODE" ]; then
    log_info "开始月度财务数据全量采集（披露季，含费用明细）"
else
    log_info "开始月度财务数据增量更新（含费用明细）"
fi
"$PYTHON_BIN" -m quantia.job.stock_financial_data $FIN_MODE --expenses >> "$LOG_FILE" 2>&1
log_info "月度财务数据更新完成"
